# app/deploy/vagrant-supervisord.conf
#############################################
#
# /etc/supervisor/conf.d/olass-supervisord.conf
#
# This file is used to maintain the
# uwsgi olass application running using supervisord.
#
# $ sudo supervisorctl restart all
# $ sudo less /var/log/supervisor/supervisord.log
# $ sudo less /var/log/uwsgi/olass.log
#
##############################################
[program:olass]
# user=supervisord
command = /usr/bin/uwsgi /srv/apps/olass/app/deploy/vagrant-uwsgi.ini

autostart=true
autorestart=true
redirect_stderr=true
startretries=1
exitcodes=0
